# supervisord.conf
[supervisord]
nodaemon=true

[program:network_watcher]
command=python3 /network_watcher.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
